# deploy/supervisor/orion-reverb.conf
[program:orion-reverb]
process_name=%(program_name)s
command=php /var/www/orion/artisan reverb:start --port=8080
directory=/var/www/orion
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/www/orion/storage/logs/reverb.log
stopwaitsecs=10
