# deployment/supervisor/orion.conf
; ============================================================
; ORION — Supervisor configuration
; Deploy to: /etc/supervisor/conf.d/orion.conf
; After deploy: supervisorctl reread && supervisorctl update
; ============================================================

; --- Laravel Reverb WebSocket Server ---
[program:orion-reverb]
command=/usr/bin/php /var/www/orion/artisan reverb:start --host=127.0.0.1 --port=8080 --no-interaction
directory=/var/www/orion
user=www-data
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stdout_logfile=/var/log/supervisor/orion-reverb.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/orion-reverb-error.log
stderr_logfile_maxbytes=5MB

; --- ChirpStack MQTT Subscriber ---
[program:orion-mqtt]
command=/usr/bin/php /var/www/orion/artisan mqtt:subscribe
directory=/var/www/orion
user=www-data
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=15
stdout_logfile=/var/log/supervisor/orion-mqtt.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/orion-mqtt-error.log
stderr_logfile_maxbytes=5MB

; --- Laravel Queue Worker (broadcasts & jobs) ---
[program:orion-queue]
command=/usr/bin/php /var/www/orion/artisan queue:work --sleep=3 --tries=3 --max-time=3600
directory=/var/www/orion
user=www-data
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=60
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/var/log/supervisor/orion-queue.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/var/log/supervisor/orion-queue-error.log
stderr_logfile_maxbytes=5MB
